# launchd/com.jurrejan.reminders-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.jurrejan.reminders-bridge</string>

	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/uv</string>
		<string>run</string>
		<string>--directory</string>
		<string>/Users/jurrejan/Documents/development/python/reminders-bridge</string>
		<string>rbridge</string>
		<string>run</string>
	</array>

	<key>KeepAlive</key>
	<true/>

	<key>RunAtLoad</key>
	<true/>

	<key>StandardOutPath</key>
	<string>/Users/jurrejan/Library/Logs/reminders-bridge.log</string>

	<key>StandardErrorPath</key>
	<string>/Users/jurrejan/Library/Logs/reminders-bridge.log</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/jurrejan</string>
		<key>PATH</key>
		<!-- ~/.local/bin holds the `bd` CLI; without it the daemon cannot
		     reconcile beads (FileNotFoundError: 'bd'). launchd does not expand
		     ~, so the path is absolute. -->
		<string>/Users/jurrejan/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
		<key>RBRIDGE_POLL_S</key>
		<string>5</string>
	</dict>
</dict>
</plist>
